# units/user/systemd-network-manager.service.in
[Unit]
Description=Bind systemd targets to network states

[Service]
Slice=background.slice
Restart=always
ExecStart=LIBEXECDIR/systemd-network-manager --user

[Install]
WantedBy=default.target
