# units/system/systemd-network-manager.service.in
[Unit]
Description=Bind systemd targets to network states
After=systemd-networkd.service

[Service]
Restart=always
ExecStart=LIBEXECDIR/systemd-network-manager

[Install]
WantedBy=systemd-networkd.service
