# units/captive-portal.target
[Unit]
Description="We have a network connection but can't reach the internet."
Conflicts=online.target offline.target
