# units/captive-portal.target
[Unit]
Description=Captive Portal
Conflicts=online.target offline.target
