# units/offline.target
[Unit]
Description=Offline
Conflicts=online.target captive-portal.target
