# units/offline.target
[Unit]
Description="The network is down."
Conflicts=online.target captive-portal.target
