# units/online.target
[Unit]
Description=Online
Conflicts=offline.target captive-portal.target
