# units/online.target
[Unit]
Description="The network is up and we have a working internet connection."
Conflicts=offline.target captive-portal.target
